# debian/usr/lib/systemd/system/monad-cruft.service
[Unit]
Description=Monad cleanup: delete old ledger/WAL/forkpoint artifacts
After=network.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/monad/env
ExecStart=/opt/monad/scripts/clear-old-artifacts.sh
DynamicUser=yes
Group=monad
SupplementaryGroups=monad-ipc
ProtectHome=read-only
NoNewPrivileges=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ReadWritePaths=/var/lib/monad

[Install]
WantedBy=multi-user.target
